# leafnode/cron.d/leafnode
MAILTO=""
* 12 * * * * news /usr/bin/texpire
*/5 * * * * news /usr/bin/fetchnews
