# mail/getmail/cron.d/getmail
*/5 * * * * geekard /usr/bin/getmail --rcfile=getmailrc --rcfile=gmailrc --rcfile=baidurc &>/dev/null
